# integration/systemd/vigilia-inference.service
# integration/systemd/vigilia-inference.service
# VIGIL-IA · CORFO Semilla Inicia 25INI-282394
#
# Corre scripts/04_deploy_inference.py como proceso nativo en el Jetson
# Orin (no en Docker), para tener acceso directo a la GPU/TensorRT sin la
# complejidad adicional del runtime NVIDIA para contenedores. Escribe en la
# misma base SQLite que consume backend/ (ver docker-compose.yml, volumen
# ../data).
#
# Instalación:
#   sudo cp vigilia-inference.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vigilia-inference.service

[Unit]
Description=VIGIL-IA - Motor de inferencia (Jetson Orin)
After=network.target
Wants=vigilia-backend.service

[Service]
Type=simple
User=vigilia
WorkingDirectory=/opt/vigilia-ia
EnvironmentFile=/opt/vigilia-ia/integration/.env
ExecStart=/usr/bin/python3 scripts/04_deploy_inference.py \
    --model ${VIGILIA_MODEL_PATH} \
    --source ${VIGILIA_RTSP_SOURCE} \
    --db ${VIGILIA_DB_PATH} \
    --gpio
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
